# stickers.service
[Unit]
Description=Sticker printer interface.

[Service]
Type=simple
ExecStart=/bin/bash //stickers.sh

[Install]
WantedBy=multi-user.target
